# deploy/ansible/roles/etphonehome-server/templates/etphonehome-mcp.service.j2
# Reach MCP Server Service
# Managed by Ansible

[Unit]
Description=Reach MCP Server
After=network.target sshd.service
Wants=network-online.target

[Service]
Type=simple
User={{ reach_user }}
Group={{ reach_group }}
EnvironmentFile={{ reach_config_dir }}/server.env
ExecStart={{ reach_home }}/venv/bin/reach-server --transport http --host ${REACH_HOST} --port ${REACH_PORT}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths={{ reach_data_dir }} {{ reach_log_dir }}

[Install]
WantedBy=multi-user.target
