# scripts/etphonehome-webhook-test.service
[Unit]
Description=Reach Webhook Test Receiver
Documentation=https://github.com/jfreed-dev/reach
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=reach
Group=reach

# Working directory
WorkingDirectory=/opt/reach

# Environment
Environment=PYTHONUNBUFFERED=1

# Start command
ExecStart=/opt/reach/venv/bin/python /opt/reach/webhook_test.py

# Restart policy
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/reach
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reach-webhook

[Install]
WantedBy=multi-user.target
