# scripts/etphonehome-mcp.service
[Unit]
Description=Reach MCP Server (HTTP/SSE)
Documentation=https://github.com/jfreed-dev/reach
After=network-online.target reach-ssh.service
Wants=network-online.target

[Service]
Type=simple
User=reach
Group=reach

# Working directory
WorkingDirectory=/opt/reach

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/reach/server.env

# Start command
ExecStart=/opt/reach/venv/bin/python -m server.mcp_server \
    --transport http \
    --host 127.0.0.1 \
    --port 8765

# Restart policy
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/reach/.reach-server
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reach-mcp

[Install]
WantedBy=multi-user.target
