# scripts/reach.service
[Unit]
Description=Reach Client
Documentation=https://github.com/jfreed-dev/reach
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Environment=HOME=/home/%i
ExecStart=/usr/local/bin/reach
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/.reach
PrivateTmp=true

[Install]
WantedBy=multi-user.target
